# readingdb/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:readingdb]
command=/bin/reading-server -d /srv/smap_archiver/rdb/ -p 4242
